# src/logger.service
[Unit]
Description=Logger server

[Service]
ExecStart=/usr/sbin/logger-server.py

[Install]
WantedBy=multi-user.target